# systemd/potd.timer
[Unit]
Description=Run potd.service every 30 minutes

[Timer]
OnCalendar=*:0/5
